# src/debian/postrm
#!/bin/bash
#--------------------------------------
# Script Name:  postrm
# Version:      1.1
# Author:       skurka@ukaachen.de, shuening@ukaachen.de, akombeiz@ukaachen.de
# Date:         30 Oct 24
# Purpose:      Post-removal script for 'aktin-notaufnahme-i2b2'.
#               Executes after files are deleted. Handles cleanup of configurations and database.
#--------------------------------------

readonly OPERATION="${1}"

source_helper_script() {
  local helper_path="/usr/share/${DPKG_MAINTSCRIPT_PACKAGE}/helper.sh"
  if [[ -f "${helper_path}" ]]; then
    . "${helper_path}"
  else
    echo "Error: helper.sh not found at ${helper_path}" >&2
    exit 1
  fi
}

deactivate_php_curl_extension() {
  echo "Disabling Apache PHP cURL extension..."
  sed -i "s/^extension=curl/;extension=curl/" /etc/php/*/apache2/php.ini || true
  echo "Restarting Apache..."
}

delete_wildfly_user() {
  echo "Deleting WildFly user..."
  if id -u wildfly >/dev/null 2>&1; then
    deluser --remove-home wildfly
    echo "WildFly user deleted."
  else
    echo "WildFly user does not exist."
  fi
}

remove_systemd_files() {
  echo "Disabling and removing systemd service files..."
  systemctl disable wildfly postgresql apache2 || true
  rm -f "/lib/systemd/system/wildfly.service" "/lib/systemd/system/apache2.service" "/lib/systemd/system/postgresql.service"
  systemctl daemon-reload
}

drop_i2b2_database() {
  source_helper_script
  connect_to_psql
  wait_for_psql_connection
  if [[ $(eval "${PSQL} -l" | grep "i2b2" | wc -l) -gt 1 ]]; then
    echo "Dropping i2b2 PostgreSQL database and user..."
    eval "${PSQL} -v ON_ERROR_STOP=1" <<EOF >/dev/null
    __I2B2_DROP_STATEMENT__
EOF
  else
    echo "i2b2 database already removed. Skipping..."
  fi
}

main() {
  set -euo pipefail
  case "$OPERATION" in
    remove)
      #  Clean up tasks after a successful removal.
      deactivate_php_curl_extension
      delete_wildfly_user
      remove_systemd_files
      ;;
    purge)
      # Run during a purge to delete all configuration files.
      drop_i2b2_database
      ;;
    abort-install)
      # Used if installation is interrupted.
      ;;
    abort-upgrade)
      # Clean up after a failed upgrade.
      ;;
    abort-remove)
      # Run if removal fails or is interrupted.
      ;;
  esac
}

main "$@"
