# src/debian/postrm
#!/bin/bash
#--------------------------------------
# Script Name:  postrm
# Version:      1.1
# Authors:      skurka@ukaachen.de, shuening@ukaachen.de, akombeiz@ukaachen.de
# Date:         30 Oct 24
# Purpose:      Cleans up configurations, user accounts, and remnants of 'aktin-notaufnahme-i2b2' post-removal, notifies users about remaining
#               services, and performs database removal on purge.
#--------------------------------------

readonly OPERATION="${1}"

deactivate_php_curl_extension() {
  echo "Deactivating PHP curl extension..."
  if grep -q "^extension=curl" /etc/php/*/apache2/php.ini; then
    sed -i "s/^extension=curl/;extension=curl/" /etc/php/*/apache2/php.ini
    systemctl restart apache2
  else
    echo "Extension already deactivated"
  fi
}

delete_wildfly_user() {
  echo "Removing WildFly user..."
  if id -u wildfly >/dev/null 2>&1; then
    deluser --remove-home wildfly
  else
    echo "User does not exist"
  fi
}

delete_wildfly_remnants() {
  echo "Cleaning up WildFly files..."
  if [[ -d /opt/wildfly ]] || [[ -f /etc/wildfly/wildfly.conf ]]; then
    rm -rf /opt/wildfly
    rm -f /etc/wildfly/wildfly.conf
  else
    echo "No WildFly files found"
  fi
}

notify_user_about_services() {
  cat <<EOF
Note: The following services remain installed:
  - Apache (apache2)
  - PostgreSQL (postgresql)

To disable them if no longer needed:
  sudo systemctl disable apache2
  sudo systemctl disable postgresql
EOF
}

ource_temp_helper() {
  echo "Loading temporary helper script..."
  local temp_helper="/tmp/${DPKG_MAINTSCRIPT_PACKAGE}_helper.sh"
  if [[ ! -f "${temp_helper}" ]]; then
    echo "Error: Helper script not found at ${temp_helper}" >&2
    return 1
  fi
  source "${temp_helper}"
}

drop_i2b2_database() {
  echo "Removing i2b2 database..."
  if ! command -v psql >/dev/null 2>&1 || ! systemctl is-active --quiet postgresql; then
    echo "PostgreSQL not available"
    return 0
  fi

  source_temp_helper
  connect_to_psql
  wait_for_psql_connection
  if [[ $(eval "${PSQL} -l" | grep "i2b2" | wc -l) -gt 1 ]]; then
    eval "${PSQL} -v ON_ERROR_STOP=1" <<EOF >/dev/null
__I2B2_DROP_STATEMENT__
EOF
  else
    echo "Database already removed"
  fi
}

remove_temp_helper() {
  echo "Cleaning up temporary helpers..."
  local temp_helper="/tmp/${DPKG_MAINTSCRIPT_PACKAGE}_helper.sh"
  if [[ -f "${temp_helper}" ]]; then
    rm -f "${temp_helper}"
  fi
}

delete_wildfly_backup_logs() {
  echo "Removing backup logs..."
  local backup_dir="/var/backups/${DPKG_MAINTSCRIPT_PACKAGE}/wildfly_logs"
  if [[ -d "${backup_dir}" ]]; then
    rm -rf "${backup_dir}"
  else
    echo "No backup logs found"
  fi
}

main() {
  set -euo pipefail
  case "$OPERATION" in
    remove)
      # Runs during package removal after files have been deleted.
      deactivate_php_curl_extension
      delete_wildfly_user
      delete_wildfly_remnants
      notify_user_about_services
      ;;
    purge)
      # Runs during a purge after 'postrm remove' to delete all leftover files.
      drop_i2b2_database
      delete_wildfly_backup_logs
      remove_temp_helper
      ;;
  esac
}

main "$@"
